# scripts/systemd/juara-gdrive-sync.timer
[Unit]
Description=Run Juara Google Drive uploader every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=juara-gdrive-sync.service

[Install]
WantedBy=timers.target
