# systemd/juara-networkpi-maintenance.service
[Unit]
Description=Juara USB export and setup-network maintenance
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/juara_networkpi_maintenance
